# deploy/auto-invest.service
[Unit]
Description=auto-invest live trading worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=auto-invest
Group=auto-invest
WorkingDirectory=/opt/auto-invest
EnvironmentFile=/opt/auto-invest/.env
# Wrapper toggles dry-run / live based on AUTO_INVEST_MODE in .env.
ExecStart=/opt/auto-invest/deploy/run-worker.sh
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/opt/auto-invest/data /opt/auto-invest/logs
ProtectHome=true

[Install]
WantedBy=multi-user.target
